# pdf2word.service
# ============================================================================
# PDF2Word Converter - systemd 服务文件
# ============================================================================
#
# 安装步骤:
#   1. sudo cp pdf2word.service /etc/systemd/system/
#   2. sudo systemctl daemon-reload
#   3. sudo systemctl enable --now pdf2word
#
# 常用命令:
#   sudo systemctl status pdf2word     # 查看状态
#   sudo systemctl restart pdf2word    # 重启服务
#   sudo systemctl stop pdf2word       # 停止服务
#   sudo journalctl -u pdf2word -f     # 查看实时日志
#
# ============================================================================

[Unit]
Description=PDF2Word Converter - 文档格式转换服务
Documentation=https://github.com/your-repo/pdf2word
After=network.target
Wants=network.target

[Service]
# 运行用户（部署时由 deploy.sh 自动替换为实际用户）
User=${APP_USER}
Group=${APP_USER}

# 工作目录（部署时由 deploy.sh 自动替换为实际路径）
WorkingDirectory=${APP_DIR}

# 启动命令
ExecStart=${APP_DIR}/venv/bin/python -m uvicorn backend.app.main:app --host 0.0.0.0 --port 8000 --workers 1

# 环境变量
Environment=PATH=${APP_DIR}/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=-${APP_DIR}/backend/.env

# 进程管理
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# 重启策略
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pdf2word

# 为 LibreOffice 提供可写的临时目录
RuntimeDirectory=pdf2word
Environment=TMPDIR=/run/pdf2word

# --- 可选安全加固（如服务能正常运行，可逐项取消注释）---
# NoNewPrivileges=true
# PrivateTmp=true
# ProtectSystem=full
# ProtectHome=read-only
# ReadWritePaths=${APP_DIR}/backend/uploads ${APP_DIR}/backend/converted /run/pdf2word
# MemoryMax=512M

[Install]
WantedBy=multi-user.target
