# pdf2word.service
# ============================================================================
# PDF2Word Converter - systemd 服务文件
# ============================================================================
#
# 安装步骤:
#   1. 替换 ${APP_DIR} 为项目实际路径（如 /opt/pdf2word）
#   2. 替换 ${APP_USER} 为运行用户（如 www-data 或当前用户名）
#   3. sudo cp pdf2word.service /etc/systemd/system/
#   4. sudo systemctl daemon-reload
#   5. sudo systemctl enable --now pdf2word
#
# 常用命令:
#   sudo systemctl status pdf2word     # 查看状态
#   sudo systemctl restart pdf2word    # 重启服务
#   sudo systemctl stop pdf2word       # 停止服务
#   sudo journalctl -u pdf2word -f     # 查看实时日志
#
# ============================================================================

[Unit]
Description=PDF2Word Converter - 文档格式转换服务
Documentation=https://github.com/your-repo/pdf2word
After=network.target
Wants=network.target

[Service]
# 运行用户和组（部署时替换为实际用户）
User=${APP_USER}
Group=${APP_USER}

# 工作目录（部署时替换为实际路径）
WorkingDirectory=${APP_DIR}

# 启动命令
ExecStart=${APP_DIR}/venv/bin/uvicorn backend.app.main:app --host 0.0.0.0 --port 8000 --workers 1

# 环境变量（可选：从文件加载）
EnvironmentFile=-${APP_DIR}/backend/.env

# 重启策略
Restart=always
RestartSec=5

# 防止无限重启循环
StartLimitInterval=60
StartLimitBurst=3

# 标准输出/错误到 journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pdf2word

# 资源限制（可选）
# MemoryMax=512M
# CPUQuota=50%

# 安全加固（可选，根据实际路径调整）
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=${APP_DIR}/backend/uploads ${APP_DIR}/backend/converted
ReadOnlyPaths=${APP_DIR}

[Install]
WantedBy=multi-user.target
